# ansible/roles/common/templates/chessus-resurrector.service.j2

[Unit]
Description=PM2 resurrect saved processes
After=network.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/bin/pm2 resurrect
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
